# Neon/backups/.config/autostart/sudo.desktop
[Desktop Entry]
Exec=sudo tlp start
Icon=system-run
Path=
Terminal=false
Type=Application
